# config/control.cron
URAVO_DIR=/opt/uravo
* * * * * uravo $URAVO_DIR/bin/control.pl&

