# config/outpost.cron
URAVO_DIR=/opt/uravo
* * * * * uravo $URAVO_DIR/bin/outpost.pl&

